# raw/ska-mid-cbf-talondx-bist/etc/systemd/system/bist.timer
[Unit]
Description=Timer for the BIST service

[Timer]
OnBootSec=60seconds

[Install]
WantedBy=timers.target