# raw/ska-mid-cbf-talondx-bist/etc/systemd/system/bist.service
[Unit]
Description=Built In Self Test (BIST) unit configuration file

[Service]
Type=oneshot
ExecStart=/bin/bist.sh -r